# services/activities-restart.service
[Unit]
Description=Restart activities.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart activities.service
